# ops/namoo-reader.service
[Unit]
Description=Namoo Reader RSS and creation site
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/rssreader
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=3088
Environment=STARTUP_REFRESH_DELAY_MS=-1
ExecStart=/usr/bin/node /opt/rssreader/server.js
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
